# supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=mongod --smallfiles
autorestart = true

[program:nodejs]
command=nodejs /yote/yote.js
autorestart = true